# .config/autostart/firstconfig.desktop
[Desktop Entry]
Name=Neovim
Exec=/usr/bin/bash -c "$HOME/.local/bin/after-reboot.sh"
Terminal=true
Type=Application
